# Library/LaunchAgents/com.diegopetrucci.mac-space-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.diegopetrucci.mac-space-cleanup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/diegopetrucci/.local/bin/scripts/mo-clean.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/diegopetrucci/Library/Logs/com.diegopetrucci.mac-space-cleanup.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/diegopetrucci/Library/Logs/com.diegopetrucci.mac-space-cleanup.err.log</string>
</dict>
</plist>
